# ansible-playbooks/monitoring/templates/alertmanager.service.j2
[Unit]
Description=AlertManager
Documentation=https://prometheus.io/docs/alerting/alertmanager/
Wants=network-online.target
After=network-online.target

[Service]
# ⚠️ 教学声明：此服务文件主要用于学习环境，生产环境请根据实际需求调整

# 用户和组
User=alertmanager
Group=alertmanager

# 工作目录
WorkingDirectory=/opt/alertmanager

# 启动命令
ExecStart=/opt/alertmanager/alertmanager \
  --config.file=/etc/alertmanager/alertmanager.yml \
  --storage.path=/var/lib/alertmanager \
  --web.listen-address=0.0.0.0:{{ alertmanager_web_port }} \
  --cluster.listen-address=0.0.0.0:{{ alertmanager_cluster_port }} \
  --web.external-url=http://localhost:{{ alertmanager_web_port }}

# 重启策略
Restart=always
RestartSec=5

# 超时设置
TimeoutStartSec=60
TimeoutStopSec=30

# 进程管理
KillMode=control-group
KillSignal=SIGTERM
SendSIGKILL=no

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/alertmanager

# 资源限制
LimitNOFILE=65536
LimitNPROC=32768

# 环境变量
Environment=HOME=/opt/alertmanager

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alertmanager

[Install]
WantedBy=multi-user.target